# devices/init/conf/rc.local
#!/bin/sh
# add your commands 
# call your scripts here

/home/ubuntu/init/init.sh

# last line must be exit 0 
exit 0
